# src/services/ss.conf
[program:ss]
command=ssserver -c "/etc/ss.json"
autostart=true
autorestart=true
stderr_logfile=/var/log/sms.err.log
stdout_logfile=/var/log/sms.out.log