# services/ms.conf
[program:ms]
command=/home/ms/movableservice/publish/MovableService
autostart=true
autorestart=true
stderr_logfile=/var/log/ms.err.log
stdout_logfile=/var/log/ms.out.log